# etc/crontab
# m h  dom mon dow   command
0 5 * * 1 ../cron.sh